# gamiscreen-client/systemd/gamiscreen-client.service
[Unit]
Description=GamiScreen Linux Client (session agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/gamiscreen-client --config %h/.config/gamiscreen/client.yaml
Restart=always
RestartSec=5s
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
